# services/ansible-cron.service
[Unit]
Description=Ansible Update at Timer
Wants=network-online.target
After=network-online.target
After=ansible-boot.service

[Service]
Type=oneshot
Environment="ANSIBLE_LOG_PATH=/root/ansible-cron.log"
ExecStartPre=/usr/bin/wget -q https://raw.githubusercontent.com/techno-link/zero-client/v3/ansible/zero.yml -O /root/zero.yml
ExecStart=/usr/bin/ansible-playbook /root/zero.yml -v

